# web/deploy/bible-memory.service
[Unit]
Description=Bible Memory Lab
After=network.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/opt/bible-memory/current
EnvironmentFile=/etc/bible-memory.env
ExecStart=/usr/bin/npm start
Restart=always
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
